# ops/cron/crontab
# Каждое ПН в 04:00 запускаем переиндексацию
0 4 * * 1 curl -fsS --retry 5 --retry-delay 5 "$API_URL"
